# deploy/systemd/openvpn-crl.path
[Unit]
Description=Watch CRL for changes and reload OpenVPN

[Path]
PathChanged=/etc/openvpn/crl.pem

[Install]
WantedBy=multi-user.target
